# deploy/systemd/run-nico-ar-nightly.service
[Unit]
Description=Nightly rebuild and deploy for run.nico.ar
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/srv/repos/personal/argensonix/labs/run-nico-ar
ExecStart=/usr/bin/bash /srv/repos/personal/argensonix/labs/run-nico-ar/scripts/deploy-nightly.sh
